# debbuilder/platform_server/DEBIAN/postrm
#!/bin/sh

if [ "$1" = "remove" ]; then
	update-rc.d platform_server remove
fi

#Restart after every install/update
service rsyslog reload > /dev/null 2>&1 ||:
