# debbuilder/platform_server/DEBIAN/prerm
#!/bin/sh

if [ "$1" = remove ]; then
  systemctl disable platform_server
  systemctl stop platform_server
fi

