# debbuilder/platform_server/DEBIAN/prerm
#!/bin/sh

if [ "$1" = remove ]; then
  invoke-rc.d platform_server stop
fi

