# docker_build/etc/supervisor/conf.d/nuboplatformserver.conf
[program:nuboplatformserver]
directory=/opt/platform_server
command=/usr/bin/node dist/restserver.js
user=root
autorestart=true
priority=400
