# docker_build/etc/supervisor/conf.d/pulseaudio.conf
[program:pulseaudio]
directory=/home/nubo
command=/usr/bin/pulseaudio
user=nubo
environment=HOME="/home/nubo",USER="nubo"
autorestart=true
priority=600
