# docker_build/etc/supervisor/conf.d/rsyslog.conf
[program:rsyslog]
command=/usr/sbin/rsyslogd -n
user=root
autorestart=true
priority=200
